# playright_v1/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ---- Xvfb: Virtueller Bildschirm ----
[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac +extension GLX +render -noreset
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

; ---- Fluxbox: Window Manager ----
[program:fluxbox]
command=fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20
startsecs=3
stdout_logfile=/var/log/fluxbox.log
stderr_logfile=/var/log/fluxbox.err

; ---- x11vnc: VNC-Server ----
[program:x11vnc]
command=x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -xkb
autorestart=true
priority=30
startsecs=5
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err

; ---- noVNC: Web-VNC-Client ----
[program:novnc]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true
priority=40
startsecs=5
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err

; ---- FastAPI: Python-API ----
[program:fastapi]
command=python3 -m uvicorn api_server:app --host 0.0.0.0 --port 8080
directory=/app
environment=DISPLAY=":99",PYTHONUNBUFFERED="1",PYTHONPATH="/app/src"
autorestart=true
priority=50
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
